# cron/backup.crontab
4 4 * * 1-5 root . /etc/default/backup && /usr/bin/backup >> /var/log/backup.log 2>&1
