# systemd/backup@.timer
[Unit]
Description=Daily full system backup

[Timer]
OnCalendar=Tue-Sat *-*-* 04:04:04
Persistent=true

[Install]
WantedBy=timers.target
